# target/linux/ramips/mt7621/base-files/etc/hotplug.d/ieee80211/10_fix_wifi_mac
#!/bin/ash

[ "$ACTION" = "add" ] || exit 0

PHYNBR=${DEVPATH##*/phy}

[ -n $PHYNBR ] || exit 0

. /lib/functions.sh
. /lib/functions/system.sh

board=$(board_name)

case "$board" in
	jdcloud,re-sp-01b)
		label_mac=$(mtd_get_mac_ascii config "mac")
		echo $(macaddr_add $label_mac $PHYNBR) > /sys${DEVPATH}/macaddress
		;;
esac
